# systemd/discord-mpris-rs.service
[Unit]
Description=Discord rich presence integration for MPRIS.
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/discord-mpris-rs
EnvironmentFile=%h/.config/discord-mpris-rs/.env

[Install]
WantedBy=default.target